# deploy/privacy-tracker.service
# systemd unit for the Django app. Install with:
#   sudo cp deploy/privacy-tracker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now privacy-tracker
[Unit]
Description=Privacy Policy Tracker (gunicorn)
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/privacy-tracker/privacy_policy_tracker/website
EnvironmentFile=/etc/privacy-tracker.env
ExecStart=/opt/privacy-tracker/privacy_policy_tracker/website/.venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 120 \
    config.wsgi:application
Restart=on-failure

[Install]
WantedBy=multi-user.target
